# etc/factoids-web.conf
start on started

script
  exec /usr/bin/java -jar /usr/local/www/gianni/b-ot.gf3.ca/factoids-web.jar >> /usr/local/www/gianni/b-ot.gf3.ca/log/factoids-web.log 2>&1
end script
